# infra/scheduled_reports.timer
[Unit]
Description=Run PesaGuard scheduled reports daily

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
